# ops/launchd/com.0luka.dispatcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.0luka.dispatcher</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python3</string>
      <string>-m</string>
      <string>core</string>
      <string>dispatch</string>
      <string>--watch</string>
      <string>--interval</string>
      <string>2</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/icmini/0luka</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>LUKA_RUNTIME_ROOT</key>
      <string>/Users/icmini/0luka_runtime</string>
      <key>ROOT</key>
      <string>/Users/icmini/0luka</string>
      <key>PYTHONUNBUFFERED</key>
      <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/icmini/0luka_runtime/logs/dispatcher.launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/icmini/0luka_runtime/logs/dispatcher.launchd.err.log</string>
  </dict>
</plist>
