# ops/decommissioned/com.02luka.mary-gateway-v3.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.02luka.mary-gateway-v3</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/icmini/0luka/agents/mary_router/gateway_v3_router.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/icmini/0luka</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/icmini/0luka/logs/mary-gateway-v3.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/icmini/0luka/logs/mary-gateway-v3.stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>LUKA_SOT</key>
        <string>/Users/icmini/0luka</string>
        <key>PYTHONPATH</key>
        <string>/Users/icmini/0luka</string>
    </dict>
</dict>
</plist>
