# launchd/com.0luka.mission-scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.0luka.mission-scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/icmini/0luka/runtime/services/mission_scheduler_runner.zsh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>LUKA_RUNTIME_ROOT</key>
        <string>/Users/icmini/0luka_runtime</string>
        <key>PYTHONPATH</key>
        <string>/Users/icmini/0luka</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>/Users/icmini/0luka</string>

    <key>StandardOutPath</key>
    <string>/Users/icmini/0luka_runtime/logs/mission_scheduler.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/icmini/0luka_runtime/logs/mission_scheduler.stderr.log</string>
</dict>
</plist>
