# launchd/com.0luka.nas-ai-workspace-mount.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.0luka.nas-ai-workspace-mount</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>mount | grep -q '/nas/ai-workspace' || /opt/homebrew/bin/dotenvx run -f /Users/icmini/0luka/.xenv -- /bin/zsh -c 'PW=$(python3 -c "import urllib.parse,os; print(urllib.parse.quote(os.environ[\"NAS_PASSWORD\"],safe=\"\"))") && mount_smbfs "//${NAS_USER}:${PW}@${NAS_IP}/${NAS_SHARED_FOLDER}" /Users/icmini/nas/ai-workspace'</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>120</integer>
    <key>StandardOutPath</key>
    <string>/Users/icmini/0luka/observability/logs/nas_mount.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/icmini/0luka/observability/logs/nas_mount.stderr.log</string>
</dict>
</plist>
