# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/ai.local.processor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.local.processor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/icmini/01luka/gc_gg_coordination/local_ai_processor.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/icmini/01luka/gc_gg_coordination</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/icmini/01luka/logs/local_ai_processor/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/icmini/01luka/logs/local_ai_processor/stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
