# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.dashboard.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>LUKA_SOT</key>
		<string>/Users/icmini/02luka</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.02luka.dashboard.server</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>cd "/Users/icmini/02luka/g" &amp;&amp; python3 -m http.server 8088</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/icmini/02luka/logs/dashboard_server.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/icmini/02luka/logs/dashboard_server.out.log</string>
	<key>ThrottleInterval</key>
	<integer>30</integer>
	<key>WorkingDirectory</key>
	<string>/Users/icmini/02luka/g</string>
</dict>
</plist>
