# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.telegram-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>REDIS_HOST</key>
		<string>127.0.0.1</string>
		<key>REDIS_PASSWORD</key>
		<string>gggclukaic</string>
		<key>REDIS_PORT</key>
		<string>6379</string>
		<key>REDIS_URL</key>
		<string>redis://:gggclukaic@127.0.0.1:6379/0</string>
		<key>TELEGRAM_BOT_TOKEN</key>
		<string>7683891739:AAFnHw3zkXBZecmzIty4iYTv8e_QxZYgAmk</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.02luka.telegram-bridge</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/python3</string>
		<string>/Users/icmini/02luka/g/tools/redis_to_telegram.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/icmini/02luka/logs/com.02luka.telegram-bridge.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/icmini/02luka/logs/com.02luka.telegram-bridge.out.log</string>
	<key>StartInterval</key>
	<integer>5</integer>
	<key>ThrottleInterval</key>
	<integer>10</integer>
	<key>WorkingDirectory</key>
	<string>/Users/icmini/02luka</string>
</dict>
</plist>
