# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.guard-health.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.02luka.guard-health.daily</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>~/02luka/g/tools/guard_health_daily.zsh</string>
    </array>
    <!-- รันทุกวัน 08:30 -->
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
    <!-- รันทันทีตอน login รอบแรกด้วย -->
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/icmini/02luka/logs/guard_health_daily.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/icmini/02luka/logs/guard_health_daily.err</string>
  </dict>
</plist>
