# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.wo_executor.codex.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.02luka.wo_executor.codex</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-o</string>
		<string>no_nomatch</string>
		<string>-lc</string>
		<string>for f in "/Users/icmini/02luka/bridge/inbox/codex"/WO-*.zsh; do
      [[ -f "$f" ]] || continue;
      "/Users/icmini/02luka/g/tools/codex_run_wo.zsh" "$f";
      mv -f "$f" "/Users/icmini/02luka/bridge/outbox/CODEX/";
    done</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/icmini/02luka/logs/codex_wo_executor.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/icmini/02luka/logs/codex_wo_executor.out.log</string>
	<key>StartInterval</key>
	<integer>120</integer>
</dict>
</plist>
