# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.opal-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>LUKA_SOT</key>
		<string>/Users/icmini/02luka</string>
		<key>PATH</key>
		<string>/Users/icmini/02luka/.venv/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>PYTHONPATH</key>
		<string>/Users/icmini/02luka:/Users/icmini/02luka/g/apps/opal_api</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.02luka.opal-api</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/icmini/02luka/.venv/bin/uvicorn</string>
		<string>opal_api_server:app</string>
		<string>--host</string>
		<string>0.0.0.0</string>
		<string>--port</string>
		<string>7001</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/icmini/02luka/logs/opal_api.stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/icmini/02luka/logs/opal_api.stdout.log</string>
	<key>ThrottleInterval</key>
	<integer>10</integer>
	<key>WorkingDirectory</key>
	<string>/Users/icmini/02luka/g/apps/opal_api</string>
</dict>
</plist>
