# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.ci-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.02luka.ci-watcher</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>source ~/02luka/paths.env &amp;&amp; source ~/02luka/g/tools/ci_watcher_config.zsh &amp;&amp; ~/02luka/g/tools/ci_watcher.sh &gt;&gt; ~/02luka/logs/ci_watcher.log 2&gt;&amp;1</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>~/02luka/logs/ci_watcher.err</string>
	<key>StandardOutPath</key>
	<string>~/02luka/logs/ci_watcher.out</string>
	<key>StartInterval</key>
	<integer>300</integer>
	<key>ThrottleInterval</key>
	<integer>60</integer>
</dict>
</plist>
