# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.distributedtask.dispatcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.distributedtask.dispatcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/icmini/01luka/distributed_task_system/start_dispatcher.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/icmini/01luka/distributed_task_system</string>
    <key>StandardOutPath</key>
    <string>/Users/icmini/01luka/logs/distributed_task/dispatcher/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/icmini/01luka/logs/distributed_task/dispatcher/stderr.log</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
