# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.memory.hub.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>LUKA_SOT</key>
		<string>/Users/icmini/02luka</string>
		<key>REDIS_HOST</key>
		<string>localhost</string>
		<key>REDIS_PASSWORD</key>
		<string>changeme-02luka</string>
		<key>REDIS_PORT</key>
		<string>6379</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.02luka.memory.hub</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>/Users/icmini/02luka/agents/memory_hub/memory_hub.py</string>
		<string>/Users/icmini/02luka/.venv/bin/python3</string>
		<string>~/02luka/agents/memory_hub/memory_hub.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/icmini/02luka/logs/memory_hub.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/icmini/02luka/logs/memory_hub.out.log</string>
	<key>ThrottleInterval</key>
	<integer>30</integer>
</dict>
</plist>
