# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.ci-coordinator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.02luka.ci-coordinator</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/env</string>
		<string>bash</string>
		<string>-lc</string>
		<string>source $HOME/02luka/paths.env &amp;&amp; node $HOME/02luka/g/tools/ci/ci_coordinator.cjs</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>$HOME/02luka/logs/ci_coordinator.err</string>
	<key>StandardOutPath</key>
	<string>$HOME/02luka/logs/ci_coordinator.out</string>
	<key>ThrottleInterval</key>
	<integer>60</integer>
</dict>
</plist>
