# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.lac-manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>LAC_BASE_DIR</key>
		<string>/Users/icmini/02luka</string>
	</dict>
	<key>Label</key>
	<string>com.02luka.lac-manager</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>/Users/icmini/02luka/agents/lac_manager/lac_manager.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/icmini/02luka/g/logs/launchd_lac_manager.err</string>
	<key>StandardOutPath</key>
	<string>/Users/icmini/02luka/g/logs/launchd_lac_manager.out</string>
	<key>StartInterval</key>
	<integer>10</integer>
</dict>
</plist>
