# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.memory.metrics.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.02luka.memory.metrics</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>export LUKA_SOT="$HOME/02luka"; "$HOME/02luka/g/tools/memory_metrics.zsh"</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>~/02luka/logs/memory_metrics.err.log</string>
	<key>StandardOutPath</key>
	<string>~/02luka/logs/memory_metrics.out.log</string>
	<key>StartInterval</key>
	<integer>3600</integer>
	<key>ThrottleInterval</key>
	<integer>30</integer>
</dict>
</plist>
