# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.phase15.quickhealth.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.02luka.phase15.quickhealth</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>TZ=Asia/Bangkok &quot;$HOME/02luka/tools/phase15_quick_health.zsh&quot; --json &gt; &quot;$HOME/02luka/mls/status/phase15_quickhealth.json&quot; 2&gt;&amp;1; exit 0</string>
  </array>
  <key>StartInterval</key>
  <integer>600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>StandardOutPath</key>
  <string>$HOME/Library/Logs/phase15_quickhealth.log</string>
  <key>StandardErrorPath</key>
  <string>$HOME/Library/Logs/phase15_quickhealth.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>TZ</key>
    <string>Asia/Bangkok</string>
  </dict>
  <key>ThrottleInterval</key>
  <integer>30</integer>
</dict>
</plist>

