# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.ggmesh.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ggmesh.watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/icmini/01luka/.venv_luka/bin/python</string>
        <string>/Users/icmini/01luka/agents/watchdog/main.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/icmini/01luka/logs/watchdog/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/icmini/01luka/logs/watchdog/stderr.log</string>
</dict>
</plist>
