# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.sync.gdrive.4h.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.02luka.sync.gdrive.4h</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>/Users/icmini/02luka/g/tools/backup_to_gdrive.zsh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/icmini/02luka/logs/backup_to_gdrive_stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/icmini/02luka/logs/backup_to_gdrive_stdout.log</string>
	<key>StartInterval</key>
	<integer>14400</integer>
</dict>
</plist>
