# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.0luka.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.0luka.heartbeat</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/icmini/0luka/.0luka/scripts/atg_multi_snap.zsh</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>StandardOutPath</key>
    <string>/Users/icmini/02luka/logs/heartbeat.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/icmini/02luka/logs/heartbeat.err</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
