# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.mls.ledger.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.02luka.mls.ledger.monitor</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/icmini/02luka/tools/mls_ledger_monitor.zsh</string>
		<string>/Users/icmini/02luka/g/tools/mls_ledger_monitor.zsh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/icmini/02luka/logs/mls_ledger_monitor.stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/icmini/02luka/logs/mls_ledger_monitor.stdout.log</string>
	<key>StartInterval</key>
	<integer>3600</integer>
	<key>ThrottleInterval</key>
	<integer>30</integer>
</dict>
</plist>
