# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.backup.gdrive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.02luka.backup.gdrive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/icmini/02luka/g/tools/backup_to_gdrive.zsh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/icmini/02luka/logs/backup_gd.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/icmini/02luka/logs/backup_gd.out.log</string>
	<key>StartInterval</key>
	<integer>28800</integer>
</dict>
</plist>
