# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.mls_watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.02luka.mls_watcher</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/icmini/02luka/g/tools/mls_file_watcher.zsh</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>LAC_BASE_DIR</key>
        <string>/Users/icmini/02luka</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    
    <key>WorkingDirectory</key>
    <string>/Users/icmini/02luka</string>
    
    <key>StandardOutPath</key>
    <string>/Users/icmini/02luka/g/logs/mls_watcher.out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/icmini/02luka/g/logs/mls_watcher.err.log</string>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
