# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.n8n.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>N8N_ENFORCE_SETTINGS_FILE_PERMISSIONS</key>
		<string>false</string>
		<key>N8N_HOST</key>
		<string>0.0.0.0</string>
		<key>N8N_PORT</key>
		<string>5678</string>
		<key>N8N_USER_FOLDER</key>
		<string>/Users/icmini/02luka/.n8n</string>
		<key>PATH</key>
		<string>/opt/homebrew/opt/node@24/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.02luka.n8n.server</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/opt/node@24/bin/node</string>
		<string>/opt/homebrew/lib/node_modules/n8n/bin/n8n</string>
		<string>start</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/icmini/02luka/logs/n8n.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/icmini/02luka/logs/n8n.out.log</string>
	<key>ThrottleInterval</key>
	<integer>30</integer>
	<key>WorkingDirectory</key>
	<string>/Users/icmini/02luka</string>
</dict>
</plist>
