# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.atg_runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.02luka.atg_runner</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>zsh</string>
    <string>/Users/icmini/02luka/tools/atg_runner_daemon.zsh</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/icmini/02luka/g/logs/atg_runner/launchd.stdout.log</string>
  <key>StandardErrorPath</key><string>/Users/icmini/02luka/g/logs/atg_runner/launchd.stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>REPO</key><string>/Users/icmini/02luka</string>
  </dict>
</dict>
</plist>
