# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.dashboard.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.02luka.dashboard.daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>$HOME/02luka/g/tools/dashboard.zsh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>0</integer>
    <key>Minute</key><integer>40
0</integer>
  </dict>
  <key>StandardOutPath</key><string>$HOME/02luka/logs/dashboard.out.log</string>
  <key>StandardErrorPath</key><string>$HOME/02luka/logs/dashboard.err.log</string>
  <key>RunAtLoad</key><true/>
</dict>
</plist>
