# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.mls-symlink-guard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.02luka.mls-symlink-guard</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/icmini/02luka/tools/guardrails/mls_symlink_guard.py</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>ROOT</key>
    <string>/Users/icmini/02luka</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/icmini/02luka_ws/g/telemetry/mls_symlink_guard.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/icmini/02luka_ws/g/telemetry/mls_symlink_guard.stderr.log</string>
</dict>
</plist>
