# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.rag.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.02luka.rag.api</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/icmini/02luka/g/rag/run_api.zsh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/icmini/02luka/logs/rag_api.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/icmini/02luka/logs/rag_api.stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>LUKA_SOT</key>
    <string>/Users/icmini/02luka</string>
    <key>RAG_HOST</key>
    <string>127.0.0.1</string>
    <key>RAG_PORT</key>
    <string>8765</string>
  </dict>
</dict>
</plist>
