# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/com.02luka.rnd.gate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.02luka.rnd.gate</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>~/02luka/g/tools/rnd_score_and_gate.zsh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>~/02luka/logs/rnd_gate.err.log</string>
	<key>StandardOutPath</key>
	<string>~/02luka/logs/rnd_gate.out.log</string>
	<key>StartInterval</key>
	<integer>420</integer>
	<key>ThrottleInterval</key>
	<integer>30</integer>
</dict>
</plist>
