# observability/artifacts/root_artifacts_legacy/backups/plist_migration_1769206548/org.02luka.sot.render.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.02luka.sot.render</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/icmini/Library/02luka_runtime/tools/sot_render.sh</string>
	</array>
	<key>StandardErrorPath</key>
	<string>/Users/icmini/Library/Logs/02luka/02luka.sot.render.err</string>
	<key>StandardOutPath</key>
	<string>/Users/icmini/Library/Logs/02luka/02luka.sot.render.out</string>
	<key>StartInterval</key>
	<integer>43200</integer>
</dict>
</plist>
